# files/etc/supervisord.d/certum.ini
[group:certumsigner]
programs=xvnc,metacity,stalonetray,ssd,p11-kit

[program:xvnc]
command=/usr/bin/Xvnc :99 -geometry 1024x768 -depth 24 -rfbport 5900 -rfbauth /root/.vnc/passwd -SecurityTypes VncAuth
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=100

[program:metacity]
command=/usr/bin/metacity
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=110

[program:stalonetray]
command=/usr/bin/stalonetray
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=120

[program:ssd]
command=/opt/SimplySignDesktop/SimplySignDesktop_start
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=130

[program:p11-kit]
command=/usr/bin/p11-kit server -v -f -n /run/p11-kit/p11kit.sock pkcs11:token=Code%%20Signing
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=140
